# systemd/cellular-dialup.service
[Unit]
Description=cns_rpi 5G链路监测与自恢复服务
After=network-pre.target
Before=network.target

[Service]
Type=simple
User=dcdw
RuntimeDirectory=cns-rpi
RuntimeDirectoryMode=0700
ExecStart=/usr/bin/python3 /home/dcdw/cns_rpi/scripts/cellular_dialup.py /var/lib/cns-rpi/config.json
Restart=on-failure
RestartSec=15

[Install]
WantedBy=multi-user.target
